# darwin/Library/LaunchAgents/com.github.coreos.coreos-vagrant.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.github.coreos.coreos-vagrant</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/sh</string>
      <string>-c</string>
      <string><![CDATA[
        export PATH="/usr/local/bin:${PATH}"
        COREOS_PATH=$(readlink "${HOME}/.coreos-vagrant")
        ! which vagrant &> /dev/null || [ ! -d ${COREOS_PATH} ] && exit 0
        COREOS_INSTANCE=$(vagrant global-status | awk -v coreos_path=${COREOS_PATH} 'match($5, coreos_path) {print $1,$4}')
        if [ -z ${COREOS_INSTANCE} ]
        then
          cd ${COREOS_PATH} && vagrant up | logger -i -t COREOS_VAGRANT
        elif [ 'running' != $(echo ${COREOS_INSTANCE} | awk '{print $2}') ]
        then
          vagrant up $(echo ${COREOS_INSTANCE} | awk '{print $1}') | logger -i -t COREOS_VAGRANT
        fi
      ]]></string>
    </array>
    <key>RunAtLoad</key>
    <true/>
  </dict>
</plist>
